# supervisord/keycloak.ini
[program:keycloak]
command=/tmp/keycloak-entrypoint.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
